# packages/vikunja-mcp/vikunja_mcp-0.5.3.tar.gz/vikunja_mcp-0.5.3/vikunja-docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vikunja]
command=/app/vikunja/vikunja
directory=/app/vikunja
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=VIKUNJA_DATABASE_TYPE="sqlite",VIKUNJA_DATABASE_PATH="/db/vikunja.db"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
# Wait for Vikunja to start before nginx
startsecs=3
