# packages/wasm-cli/wasm_cli-0.14.0.tar.gz/wasm_cli-0.14.0/src/wasm/templates/systemd/app.service.j2
# Systemd service for {{ name }}
# Generated by WASM

[Unit]
Description={{ description }}
Documentation=https://github.com/Perkybeet/wasm
After=network.target

[Service]
Type=simple
User={{ user }}
Group={{ group | default(user) }}
WorkingDirectory={{ working_directory }}

# Environment
{% for key, value in environment.items() %}
Environment="{{ key }}={{ value }}"
{% endfor %}

# Command
ExecStart={{ command }}

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ name }}

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
