# packages/raphson-music-headless/raphson_music_headless-1.5.0.tar.gz/raphson_music_headless-1.5.0/raphson-music-headless.service
[Unit]
Description=Headless music player
After=network-online.target

[Service]
ExecStart=raphson-music-headless

[Install]
WantedBy=multi-user.target
