# packages/obsidian-kb/obsidian_kb-2.0.7.1.tar.gz/obsidian_kb-2.0.7.1/scripts/com.obsidian-kb.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obsidian-kb</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/.local/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/path/to/obsidian-kb</string>
        <string>python</string>
        <string>-m</string>
        <string>obsidian_kb.mcp_server</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/path/to/obsidian-kb</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/tmp/obsidian-kb.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/obsidian-kb.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/USERNAME/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>OBSIDIAN_KB_DB_PATH</key>
        <string>/Users/USERNAME/.obsidian-kb/lancedb</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>

