# packages/plexus-agent/plexus_agent-0.3.11.tar.gz/plexus_agent-0.3.11/scripts/plexus.service
[Unit]
Description=Plexus Agent
Documentation=https://docs.plexus.company
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
ExecStart=/usr/local/bin/plexus run
Restart=always
RestartSec=10
Environment=PLEXUS_CONFIG_DIR=/home/pi/.plexus

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plexus

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/.plexus

[Install]
WantedBy=multi-user.target
