# packages/mongo-charms-single-kernel/mongo_charms_single_kernel-1.8.23.tar.gz/mongo_charms_single_kernel-1.8.23/single_kernel_mongo/templates/enable-transparent-huge-pages.service.j2
[Unit]
Description=Enable Transparent Hugepages (THP)
DefaultDependencies=no
After=sysinit.target local-fs.target
Before={{service_file}}

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'echo 0 | tee /sys/kernel/mm/transparent_hugepage/khugepaged/max_ptes_none > /dev/null'
ExecStartPre=/bin/sh -c 'echo defer+madvise | tee /sys/kernel/mm/transparent_hugepage/defrag > /dev/null '
ExecStart=/bin/sh -c 'echo always | tee /sys/kernel/mm/transparent_hugepage/enabled > /dev/null'

[Install]
WantedBy=basic.target
